# roles/csrgenerator/templates/csrgenerator.socket.j2
[Unit]
Description=csrgenerator gunicorn socket

[Socket]
ListenStream={{ csrgen_listen }}

[Install]
WantedBy=sockets.target
